# packaging/debian/postinst
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-or-later
set -e

if [ "$1" = configure ]; then
    systemctl daemon-reload || true
    # migrate off a previous install.sh deployment (units in /etc, bins in
    # /usr/local) so the packaged copies are the ones that run
    if [ -f /etc/systemd/system/phanspeed.service ]; then
        systemctl disable --now phanspeed.service 2>/dev/null || true
        rm -f /etc/systemd/system/phanspeed.service \
              /etc/systemd/system/phanspeed-healthcheck.service \
              /etc/systemd/system/phanspeed-healthcheck.timer
        rm -f /usr/local/bin/phanspeedd /usr/local/bin/phanspeed \
              /usr/local/bin/phanspeed-healthcheck \
              /usr/local/bin/phanspeed-tune /usr/local/bin/phanspeed-update
        systemctl daemon-reload || true
    fi
    # UUID migration phanspeed@local -> phanspeed@asuramaya: the old directory
    # is an orphan once packaging installs under the new path (dpkg has no
    # notion that the two paths are "the same" extension). Package-owned
    # cleanup only; per-user enable/disable state is the operator's own call
    # -- see packaging/activate-uuid-migration.sh.
    if [ -d /usr/share/gnome-shell/extensions/phanspeed@local ]; then
        rm -rf /usr/share/gnome-shell/extensions/phanspeed@local
    fi
    systemctl enable --now phanspeed.service || true
    systemctl enable --now phanspeed-healthcheck.timer || true
    systemctl enable --now phanspeed-update.timer || true
    # on upgrade, --now won't reload an already-running daemon, so the new code
    # would sit inert until reboot — restart it so updates actually take effect
    systemctl try-restart phanspeed.service || true
    echo "PhanSpeed installed. Enable the pill for your user with:"
    echo "    gnome-extensions enable phanspeed@asuramaya"
    echo "  (then log out/in once on Wayland), then auto-tune with:"
    echo "    sudo phanspeed-tune --target both --apply"
fi

exit 0
