# systemd/phanspeed-healthcheck.service
[Unit]
Description=PhanSpeed healthcheck (restart the daemon if it wedged)
After=phanspeed.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/phanspeed-healthcheck
