# systemd/phanspeed.service
[Unit]
Description=PhanSpeed — Dell thermal/fan control daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/phanspeedd
# Restore a neutral thermal profile + clear status on stop/restart:
ExecStopPost=/bin/sh -c 'echo balanced > /sys/firmware/acpi/platform_profile 2>/dev/null; rm -f /run/phanspeed/status.json; true'
Restart=on-failure
RestartSec=2
RuntimeDirectory=phanspeed
RuntimeDirectoryMode=0755

# --- sandboxing ---------------------------------------------------------
# Two capabilities are kept:
#   CAP_CHOWN     — hand ownership of the control socket + status file to the
#                   logged-in user so the unprivileged pill can read status and
#                   issue (SO_PEERCRED-gated) commands.
#   CAP_SYS_RAWIO — open /dev/mem to read/write the chipset MMIO RAPL register.
#                   Dell hides a low package power cap there that the MSR view
#                   cannot see; lifting it is the only way the Perf mission can
#                   actually unleash the CPU. The daemon only ever touches the
#                   PACKAGE_RAPL_LIMIT field at the MCHBAR — /dev/mem is otherwise
#                   unused. A deliberate, scoped expansion of privilege.
# Everything else is dropped — no CAP_DAC_OVERRIDE etc.; the profile/RAPL sysfs
# nodes are root-writable by their own mode bits.
CapabilityBoundingSet=CAP_CHOWN CAP_SYS_RAWIO
AmbientCapabilities=CAP_CHOWN CAP_SYS_RAWIO
NoNewPrivileges=yes
# Read-only filesystem; strict leaves /dev /proc /sys writable, so all the
# daemon's levers under /sys still work: platform_profile, intel-rapl PL1, the
# panel/keyboard backlight and the dGPU power/control node (Endure mission).
# Config dir re-allowed below.
ProtectSystem=strict
ReadWritePaths=/etc/phanspeed
ProtectHome=tmpfs
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectHostname=yes
ProtectProc=invisible
# NOTE: ProcSubset=pid is intentionally NOT set — it hides /proc/stat, which the
# daemon reads (CPU busy%) to detect a power/PROCHOT frequency clamp. Other
# processes stay hidden via ProtectProc=invisible.
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
# No network of any kind (daemon only binds an AF_UNIX socket):
IPAddressDeny=any
# No device access except the NVIDIA nodes nvidia-smi needs (GPU power control)
# and /dev/mem for the chipset MMIO RAPL register (hidden package power cap):
DevicePolicy=closed
DeviceAllow=/dev/mem rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidia-uvm rw
KeyringMode=private
UMask=0077
# NOTE: do NOT set ProtectKernelTunables=yes — it remounts /sys read-only and
# would break writing /sys/firmware/acpi/platform_profile.

[Install]
WantedBy=multi-user.target
